# linux/gnome/autostart/imwheel.desktop
[Desktop Entry]
Type=Application
Name=Imwheel
NoDisplay=true
Exec=/usr/sbin/imwheel
