# linux/gnome/autostart/xset.desktop
[Desktop Entry]
Type=Application
Exec=xset r rate 248 48
X-GNOME-Autostart-enabled=true
Name=Typematic delay and rate
